# crontab
30	8	* * * root /usr/bin/php /scripts/run.php
