# crontab
30	8	* * * /usr/bin/php /scripts/run.php