# crontab
00	23	* * * root /usr/bin/php /scripts/run.php
